# cec-wake.service
[Unit]
Description=HDMI-CEC Wake Controller
After=graphical-session.target
ConditionPathExists=/dev/cec

[Service]
Type=notify
Environment=WAYLAND_DISPLAY=wayland-0
EnvironmentFile=-%h/.config/cec-wake/cec-wake.conf
# strip .py following standard behavior
ExecStart=/usr/libexec/cec_wake_daemon
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
NotifyAccess=main
WatchdogSec=20

[Install]
WantedBy=default.target
